# src/etc/s6-overlay/s6-rc.d/init-devcoder-mods-apt-installer/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

PKG_SCRIPT_VER="20260119"


# Source the echo module
source /etc/s6-overlay/s6-rc.d/echo_module.sh
set_module_name "init-devcoder-mods-apt-installer"
say "▶️  Package install script version ${PKG_SCRIPT_VER}"

if ! grep -Eq '^(ID|ID_LIKE)=.*(debian|ubuntu)' /etc/os-release; then
    say "Non-Debian/Ubuntu OS detected, skipping installer."
    exit 0
fi

if [[ ! -x /usr/bin/apt-get ]]; then
    say "apt-get not available, skipping installer."
    exit 0
fi

normalize_packages() {
    tr ',\n' '  ' <<< "$1"
}

install_apt_packages() {
    local SOURCE=$1
    shift
    local PACKAGES=("$@")

    if [[ ${#PACKAGES[@]} -ne 0 ]] && [[ ${PACKAGES[*]} != "" ]]; then
        say "Installing apt packages from ${SOURCE}"
        say "Packages: ${PACKAGES[*]}"
        export DEBIAN_FRONTEND="noninteractive"
        apt-get update
        if apt-get install -y --no-install-recommends \
            "${PACKAGES[@]}"; then
            say "✅ Install status: success -> ${PACKAGES[*]}"
        else
say "❌ Install status: failed -> ${PACKAGES[*]}"
        fi
    else
        say "No packages defined in ${SOURCE}, skipping."
    fi
}

install_apt_packages_from_file() {
    local FILE_PATH=$1
    IFS=' ' read -ra FILE_PACKAGES <<< "$(tr '\n' ' ' < "${FILE_PATH}")"
    rm -f "${FILE_PATH}"
    install_apt_packages "${FILE_PATH}" "${FILE_PACKAGES[@]}"
}

listen_for_devcoder_apt_files() {
    local idle_seconds=0
    shopt -s nullglob
    while (( idle_seconds < 60 )); do
        local devcoder_files=(/devcoder-apt-*)
        if (( ${#devcoder_files[@]} > 0 )); then
            for devcoder_file in "${devcoder_files[@]}"; do
                install_apt_packages_from_file "${devcoder_file}"
            done
            idle_seconds=0
        else
            sleep 1
            ((idle_seconds++))
        fi
    done
    say "No /devcoder-apt-* files for 60 seconds, exiting listener."
}

listen_for_devcoder_apt_files &

if [[ -n "${APT_PACKAGES:-}" ]]; then
    IFS=' ' read -ra ENV_PACKAGES <<< "$(normalize_packages "${APT_PACKAGES}")"
    install_apt_packages "APT_PACKAGES" "${ENV_PACKAGES[@]}"
fi
