# src/etc/s6-overlay/s6-rc.d/init-devcoder-mods/up
# This file doesn't do anything, it's just the start of the devcoder-mods init process
