# src/etc/s6-overlay/s6-rc.d/svc-devcoder-openclaw/up
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Optional: wait for other services to be ready
# For example, if openclaw needs network or other dependencies
