# src/etc/s6-overlay/s6-rc.d/init-devcoder-nodejs/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
source /etc/s6-overlay/s6-rc.d/echo_module.sh
set_module_name "init-devcoder-nodejs"



say "â–¶ï¸  Start nodejs installer"

export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:$PATH"
ENV_INSTALL_NPM_PACKAGES=""
FILE_INSTALL_NPM_PACKAGES=""
INSTALL_NPM_PACKAGES_TO_INSTALL=""

run_as_root() {
  if [ "$(id -u)" = "0" ]; then
    "$@"
    return
  fi
  if command -v sudo >/dev/null 2>&1; then
    /usr/bin/sudo -E "$@"
    return
  fi
  if command -v su >/dev/null 2>&1; then
    su -s /bin/bash root -c "$(printf '%q ' "$@")"
    return
  fi
  "$@"
}

install_nodejs_dev_environment() {
  say "**** installing nodejs dev environment ****"
  NODEJS_MOD_VERSION=${NODEJS_MOD_VERSION:-24}

  # Clean up old NodeSource repository and keyring files
  run_as_root rm -f /etc/apt/sources.list.d/nodesource.list
  run_as_root rm -f /usr/share/keyrings/nodesource.gpg
  run_as_root rm -f /usr/share/keyrings/nodesource-repo.gpg
  run_as_root rm -f /etc/apt/trusted.gpg.d/nodesource.gpg

  # Ensure keyrings directory exists
  run_as_root mkdir -p /usr/share/keyrings

  # Download and import the NodeSource signing key using official method
  run_as_root bash -c "curl -fsSL https://deb.nodesource.com/gpgkey/nodesource-repo.gpg.key | gpg --dearmor -o /usr/share/keyrings/nodesource.gpg"
  run_as_root chmod 644 /usr/share/keyrings/nodesource.gpg

  # Validate architecture
  arch=$(dpkg --print-architecture)
  if [ "$arch" != "amd64" ] && [ "$arch" != "arm64" ]; then
    say "âš ï¸  Unsupported architecture: $arch. Only amd64, arm64 are supported by NodeSource."
  fi

  # Add repository using official format
  run_as_root bash -c "echo 'deb [arch=$arch signed-by=/usr/share/keyrings/nodesource.gpg] https://deb.nodesource.com/node_${NODEJS_MOD_VERSION}.x nodistro main' > /etc/apt/sources.list.d/nodesource.list"

  # Set up APT preferences for proper package priority (from official script)
  run_as_root mkdir -p /etc/apt/preferences.d
  run_as_root bash -c "cat > /etc/apt/preferences.d/nodejs <<'EOF'
Package: nodejs
Pin: origin deb.nodesource.com
Pin-Priority: 600
EOF"

  # Add Yarn repository
  if [[ ! -f "/etc/apt/sources.list.d/yarn.list" ]]; then
    run_as_root bash -c "curl -fsSL https://dl.yarnpkg.com/debian/pubkey.gpg | gpg --dearmor | tee /usr/share/keyrings/yarnsource.gpg >/dev/null"
    run_as_root bash -c "echo 'deb [signed-by=/usr/share/keyrings/yarnsource.gpg] https://dl.yarnpkg.com/debian/ stable main' > /etc/apt/sources.list.d/yarn.list"
  fi

  say "**** Adding nodejs and yarn to package install list ****"
  cat <<'EOF' >> /devcoder-apt-nodejs.txt
nodejs
yarn
EOF
}

wait_for_npm_installer() {
  local waited=0
  while ! command -v npm >/dev/null 2>&1; do
    if (( waited >= 60 )); then
      say "âš ï¸  npm not available after 60 seconds."
      return 1
    fi
    sleep 1
    ((waited++))
  done
}

install_nvm() {
  local target_home=$1
  local nvm_dir="$target_home/.nvm"
  say "**** installing nvm ****"
  run_as_user bash -c "export HOME='$target_home'; export NVM_DIR='$nvm_dir'; if [ ! -s '$nvm_dir/nvm.sh' ]; then curl -o- https://raw.githubusercontent.com/nvm-sh/nvm/master/install.sh | bash; echo '**** loading nvm ****'; source '$nvm_dir/nvm.sh'; echo '**** installing the latest release ****'; nvm install node; else echo '**** nvm already installed, skipping ****'; fi"
}

run_as_user() {
  if id -u "$TARGET_USER" >/dev/null 2>&1; then
    if [ "$(id -u)" = "$(id -u "$TARGET_USER")" ]; then
      "$@"
      return
    fi
    if command -v sudo >/dev/null 2>&1; then
      /usr/bin/sudo -E -u "$TARGET_USER" env PATH="$PATH" "$@"
      return
    fi
    if command -v su >/dev/null 2>&1; then
      su -s /bin/bash "$TARGET_USER" -c "$(printf '%q ' "$@")"
      return
    fi
    if command -v s6-setuidgid >/dev/null 2>&1; then
      s6-setuidgid "$TARGET_USER" "$@"
      return
    fi
  fi
  "$@"
}

if [ -n "${INSTALL_NPM_PACKAGES:-}" ]; then
  ENV_INSTALL_NPM_PACKAGES="${INSTALL_NPM_PACKAGES//\"/}"
  ENV_INSTALL_NPM_PACKAGES="${ENV_INSTALL_NPM_PACKAGES//,/ }"
  INSTALL_NPM_PACKAGES_TO_INSTALL="$INSTALL_NPM_PACKAGES_TO_INSTALL $ENV_INSTALL_NPM_PACKAGES"
fi

TARGET_USER="${CUSTOM_USER:-abc}"
if id "$TARGET_USER" >/dev/null 2>&1; then
  TARGET_UID=$(id -u "$TARGET_USER")
  TARGET_GID=$(id -g "$TARGET_USER")
  TARGET_HOME=$(getent passwd "$TARGET_USER" | cut -d: -f6)
else
  TARGET_UID=911
  TARGET_GID=911
  TARGET_HOME="$HOME"
fi

NPM_CONFIG_FILE="/config/npm.txt"
if [ -f "$NPM_CONFIG_FILE" ]; then
  FILE_INSTALL_NPM_PACKAGES=$(parse_config_file "$NPM_CONFIG_FILE")
  INSTALL_NPM_PACKAGES_TO_INSTALL="$INSTALL_NPM_PACKAGES_TO_INSTALL $FILE_INSTALL_NPM_PACKAGES"
fi

if [ -z "$INSTALL_NPM_PACKAGES_TO_INSTALL" ]; then
  say "â—¾ No INSTALL_NPM_PACKAGES configured; skipping."
  exit 0
fi

install_nodejs_dev_environment
wait_for_npm_installer || exit 0
install_nvm "$TARGET_HOME"

INSTALL_NPM_PACKAGES_TO_INSTALL=$(printf "%s" "$INSTALL_NPM_PACKAGES_TO_INSTALL" | tr ',\n\r\t' '    ' | xargs)
say "ðŸ“¦ Installing npm packages: $INSTALL_NPM_PACKAGES_TO_INSTALL"

run_as_root chown -R "$TARGET_UID:$TARGET_GID" "$TARGET_HOME/.nvm" "$TARGET_HOME/.npm" 2>/dev/null || \
  say "âš ï¸  warn: unable to chown $TARGET_HOME/.nvm or $TARGET_HOME/.npm (continuing)" >&2

export NVM_DIR="$TARGET_HOME/.nvm"
if [ -s "$NVM_DIR/nvm.sh" ]; then
  source "$NVM_DIR/nvm.sh"
fi

if ! command -v npm >/dev/null 2>&1; then
  NPM_BIN=$(find "$NVM_DIR/versions/node" -maxdepth 8 -type f -name npm 2>/dev/null | head -1 || true)
  if [ -n "${NPM_BIN:-}" ]; then
    NODE_DIR=$(dirname "$(dirname "$NPM_BIN")")
    export PATH="$NODE_DIR/bin:$PATH"
    export NPM_CONFIG_PREFIX="$NODE_DIR"
  fi
fi

if ! command -v npm >/dev/null 2>&1; then
  say "âš ï¸  npm not found; skipping npm installs."
  exit 0
fi

if [ -n "$NVM_DIR" ] && [ -d "$NVM_DIR/versions/node" ]; then
  NODE_VERSION=$(basename "$(find "$NVM_DIR/versions/node" -maxdepth 1 -type d ! -name node | head -1)" 2>/dev/null || echo "")
  if [ -n "$NODE_VERSION" ]; then
    export NPM_CONFIG_PREFIX="$NVM_DIR/versions/node/$NODE_VERSION"
  fi
fi


for pkg in $INSTALL_NPM_PACKAGES_TO_INSTALL; do
  if run_as_user npm list -g "$pkg" >/dev/null 2>&1; then
    say "$pkg is already installed, skipping..."
  else
    say "npm install -g $pkg"
    if ! run_as_user npm install -g "$pkg"; then
      say "âš ï¸  warn: failed to install $pkg (continuing)" >&2
    fi
  fi
done
